# docker/supervisord.conf
[supervisord]
nodaemon=true

[program:apache]
command=/usr/bin/http.sh -c /etc/http.sh.conf
